# systemd/openqa-worker-cacheservice.service
[Unit]
Description=OpenQA Worker Cache Service
After=network-online.target
Wants=network-online.target

[Service]
Restart=on-failure
User=_openqa-worker
ExecStart=/usr/share/openqa/script/openqa-workercache daemon -m production

[Install]
WantedBy=multi-user.target
