# systemd/openqa-minion-restart.path
[Path]
PathChanged=installvendorlib/Minion/Backend/resources/migrations/pg.sql

[Install]
WantedBy=multi-user.target
