# systemd/openqa-worker.target
[Unit]
Description=openQA Worker
Wants=var-lib-openqa-share.mount

[Install]
WantedBy=multi-user.target
