# systemd/openqa-enqueue-needle-ref-cleanup.timer
[Unit]
Description=Enqueues a needle refs task for openQA every day.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
