# systemd/openqa-enqueue-bug-cleanup.timer
[Unit]
Description=Enqueues a bug cleanup task for the openQA database every week.

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
