# systemd/openqa-worker.target
[Unit]
Description=openQA Worker

[Install]
WantedBy=multi-user.target
