# systemd/openqa-reload-worker-auto-restart@.path
# see https://www.freedesktop.org/software/systemd/man/systemd.path.html#PathExists=

[Path]
PathChanged=/etc/openqa/workers.ini

[Install]
WantedBy=multi-user.target
