# systemd/openqa-enqueue-result-cleanup.timer
[Unit]
Description=Enqueues a result/screenshots/logs cleanup task for openQA every hour.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
